# etc/systemd/system/apple-sysrq.service
[Unit]
After=

[Service]
ExecStart=/usr/local/sbin/apple-sysrq

[Install]
WantedBy=default.target
