# etc/systemd/system/apple-backlight.service
[Unit]
After=

[Service]
ExecStart=/usr/local/sbin/backlight-allow-user

[Install]
WantedBy=default.target
